# templates/default/init.server.erb
#!/bin/sh -e
#
# rc file for TeamCity Server
#
# chkconfig: 345 96 10
# description: TeamCity Server (http://jetbrains.com/teamcity/)
#
### BEGIN INIT INFO
# Provides: teamcity
# Required-Start: $network
# Required-Stop: $network
# Default-Start: 3 4 5
# Default-Stop: 0 1 2 6
# Short-Description: TeamCity Server (http://jetbrains.com/teamcity/)
# Description: TeamCity Server (http://jetbrains.com/teamcity/)
### END INIT INFO

# Source networking configuration.
if [ -f /etc/sysconfig/teamcity-server ] ; then
  . /etc/sysconfig/teamcity-server
else
  echo "No teamcity server configuration"
  exit 0
fi

case "$1" in
  start)
        su -s /bin/bash -p -c "cd $TEAMCITY_DIR && bin/teamcity-server.sh start" $TEAMCITY_USER
        ;;
  stop)
        su -s /bin/bash -p -c "cd $TEAMCITY_DIR && bin/teamcity-server.sh stop" $TEAMCITY_USER
        ;;
  restart)
        stop
        start
        ;;
  *)
        echo $"Usage: $0 {start|stop|restart}"
        exit 1
esac
        
exit 0