# templates/default/init.server.erb
#!/bin/sh -e
#
# rc file for TeamCity Server
#
# chkconfig: 345 96 10
# description: TeamCity Server (http://jetbrains.com/teamcity/)
#
### BEGIN INIT INFO
# Provides: teamcity
# Required-Start: $network
# Required-Stop: $network
# Default-Start: 3 4 5
# Default-Stop: 0 1 2 6
# Short-Description: TeamCity Server (http://jetbrains.com/teamcity/)
# Description: TeamCity Server (http://jetbrains.com/teamcity/)
### END INIT INFO

export TEAMCITY_DATA_PATH=<%= @data_path %>
export CATALINA_PID=<%= @pid_file %>

cd <%= @path %>

case $1 in
  start)
    su -l -p -c '<%= @path %>/bin/teamcity-server.sh start' <%= @user %>
  ;;
  stop)
    su -l -p -c '<%= @path %>/bin/teamcity-server.sh stop' <%= @user %>
  ;;
  restart)
    $0 stop && $0 start
  ;;
esac

exit 0
